# themes/luci-theme-netgearmod/root/etc/uci-defaults/30_luci-theme-netgearmod
#!/bin/sh
uci batch <<-EOF
	set luci.themes.NetgearMod=/luci-static/netgearmod
	set luci.main.mediaurlbase=/luci-static/netgearmod
	commit luci
EOF
exit 0 
# ↑↑↑ DO NOT EDIT THIS LINE
